# rkn-blocker/systemd/rkn-blocker.timer
[Unit]
Description=RKN Blocklist Daily Update at 4 AM
Documentation=https://github.com/HolyLittleGirl/3xui-shopbot/tree/feature/rkn-module/rkn-blocker
Requires=rkn-blocker.service
After=network.target

[Timer]
# Первый запуск через 5 минут после старта
OnBootSec=5min
# Ежедневно в 4:00 утра
OnCalendar=*-*-* 04:00:00
# Сохранять состояние между перезагрузками
Persistent=true
Unit=rkn-blocker.service

[Install]
WantedBy=timers.target
