# rkn-blocker-inside/systemd/rkn-blocker.service
[Unit]
Description=RKN IP Blocker - Блокировка запрещённых ресурсов РФ
Documentation=https://github.com/HolyLittleGirl/3xui-shopbot/tree/feature/rkn-module/rkn-blocker
After=network.target
Wants=network.target

[Service]
Type=oneshot
ExecStart=/opt/rkn-blocker/venv/bin/python3 /opt/rkn-blocker/block_ips.py run
WorkingDirectory=/opt/rkn-blocker
StandardOutput=journal
StandardError=journal
RemainAfterExit=yes

# Security hardening
NoNewPrivileges=yes
ProtectSystem=strict
ReadWritePaths=/var/log/rkn-blocker /opt/rkn-blocker

[Install]
WantedBy=multi-user.target
